# docker_env/cnxt-supervisor.conf
[program:cnxt]
directory=/app
command=/bin/sh -c " yarn start "
